# iac/common/templates/chttp-service.service.j2
# Systemd service file for {{ item.name }}
# Generated by Ansible for {{ ansible_hostname }}

[Unit]
Description={{ item.description | default('CHTTP Static Analysis Service') }}
Documentation=https://github.com/iw2rmb/ploy/tree/main/chttp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ chttp_user }}
Group={{ chttp_user }}
WorkingDirectory={{ chttp_home }}

# Service execution
ExecStart={{ chttp_home }}/bin/{{ item.binary }} --config {{ chttp_config_dir }}/{{ item.config }}
ExecReload=/bin/kill -HUP $MAINPID

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ chttp_log_dir }} /tmp/chttp
PrivateTmp=false
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M
TasksMax=100

# Process management
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3
TimeoutStartSec=30
TimeoutStopSec=10

# Environment variables
Environment=CHTTP_CONFIG_PATH={{ chttp_config_dir }}/{{ item.config }}
Environment=CHTTP_SERVICE_NAME={{ item.name }}
Environment=CHTTP_LOG_LEVEL={{ chttp_log_level | default('info') }}
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/usr/local/lib/python3/dist-packages
Environment=HOME={{ chttp_home }}

# Logging
StandardOutput=append:{{ chttp_log_dir }}/{{ item.name }}.log
StandardError=append:{{ chttp_log_dir }}/{{ item.name }}.error.log
SyslogIdentifier={{ item.name }}

[Install]
WantedBy=multi-user.target